# crontab
0 * * * * root /usr/local/bin/python /app/webscraper.py >> /var/log/cron.log 2>&1
